# etc/s6-overlay/s6-rc.d/nginx-init/run
#!/command/execlineb -P
ln -s ../sites-available/api.growing-relic.arpa.net /etc/nginx/sites-enabled/
ln -s ../sites-available/fe.growing-relic.arpa.net /etc/nginx/sites-enabled/
