# etc/s6-overlay/s6-rc.d/init-hosts/up
#!/command/execlineb -P

redirfd -a 1 /etc/hosts
echo "127.0.1.1\tapi.growing-relic.home.arpa\n127.0.1.2\tfe.growing-relic.home.arpa"
