# etc/s6-overlay/s6-rc.d/init-nginx/up
#!/command/execlineb -P

foreground {
    s6-ln -fns ../sites-available/api.growing-relic.home.arpa /etc/nginx/sites-enabled/
}
s6-ln -fns ../sites-available/fe.growing-relic.home.arpa /etc/nginx/sites-enabled/
