# etc/s6-overlay/s6-rc.d/growing-relic/run
#!/command/execlineb -P
define domain api.growing-relic.home.arpa
define sockdir /srv/${domain}/run

fdmove -c 2 1
if { mkdir -p -m 0755 ${sockdir} }
cd /root/growing-relic
export RAILS_DEVELOPMENT_HOSTS ${domain}
/root/.rbenv/shims/bundle exec puma -b unix://${sockdir}/puma.sock
